# hypr/scripts/battery-check.service
[Unit]
Description=Battery Monitor

[Service]
ExecStart=~/.config/hypr/scripts/battery_check.sh
Restart=always

[Install]
WantedBy=default.target
